# etc/udev/rules.d/regulatory.rules
KERNEL=="regulatory*"", ACTION=="change", subsystem=="PLATFORM", RUN+="/usr/sbin/crda"
